# pi-sniffer.service
[Unit]
Description=Pi Sniffer
After=network.target

[Service]
ExecStart=/home/ian/pi-sniffer/scan 192.168.0.120 1883
WorkingDirectory=/home/ian/pi-sniffer
StandardOutput=inherit
StandardError=inherit
StartLimitIntervalSec=0
Restart=always
MemoryAccounting=true
MemoryHigh=4096K
MemoryMax=8192K

[Install]
WantedBy=multi-user.target
