# pi-sniffer.service
[Unit]
Description=Pi Sniffer
After=network.target

[Service]
ExecStart=/home/ian/pi-sniffer/scan 192.168.0.52:1883 BLF 7778
WorkingDirectory=/home/ian/pi-sniffer
StandardOutput=inherit
StandardError=inherit
Restart=always
RestartSec=20
MemoryAccounting=true
MemoryHigh=8096K
MemoryMax=12192K

[Install]
WantedBy=multi-user.target
